# phusion-baseimage/jenkins/rc.local
#!/bin/bash

chown -R jenkins:tomcat /home/jenkins
chmod -R g+w /home/jenkins
/etc/init.d/jenkins start &

exit 0
