# phusion-baseimage/mcollective-server/rc.local
#!/bin/bash

sed -e "s/##ACTIVEMQ_PORT_61613_TCP_PORT##/${ACTIVEMQ_PORT_61613_TCP_PORT}/g" /etc/mcollective/server.cfg.tpl > /etc/mcollective/server.cfg

/etc/init.d/mcollective start &

exit 0
