# phusion-baseimage/mcollective-server/rc.local
#!/bin/bash

sed -e "s/##ACTIVEMQ_SERVER/${ACTIVEMQ_SERVER}/g" /etc/mcollective/server.cfg.tpl > /etc/mcollective/server.cfg

/etc/init.d/mcollective start &

exit 0
