# tomcat/rc.local
#!/bin/bash

/etc/init.d/tomcat start &

exit 0
