# phusion-baseimage/mariadb/rc.local
#!/bin/bash

/etc/init.d/mysql start

exit 0
