# phusion-baseimage/mediatomb/rc.local
#!/bin/bash

ps -ef | grep mediatomb | grep -v grep > /dev/null

while [ $? -ne 0 ]; do
    echo "(Re)starting mediatomb ..."
    service mediatomb restart
    sleep 5s
    ps -ef | grep mediatomb | grep -v grep > /dev/null
done
